# ops/mac-mini/launchd/com.agentpay.bill.prediction-report.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.agentpay.bill.prediction-report</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/baskar_viji/hedge/ops/mac-mini/bin/bill-prediction-report-scheduled</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/baskar_viji/hedge</string>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>10</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>10</integer>
      </dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/baskar_viji/hedge/.rumbling-hedge/logs/launchd-prediction-report.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/baskar_viji/hedge/.rumbling-hedge/logs/launchd-prediction-report.err.log</string>
  </dict>
</plist>
