# ops/mac-mini/launchd/com.agentpay.bill.research-collect.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.agentpay.bill.research-collect</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/baskar_viji/hedge/ops/mac-mini/bin/bill-research-collect-scheduled</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/baskar_viji/hedge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/baskar_viji/hedge/.rumbling-hedge/logs/launchd-research-collect.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/baskar_viji/hedge/.rumbling-hedge/logs/launchd-research-collect.err.log</string>
  </dict>
</plist>
