# ops/mac-mini/launchd/com.agentpay.bill.prediction-collect.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.agentpay.bill.prediction-collect</string>
    <key>ProgramArguments</key>
    <array>
      <string>__BILL_ROOT__/ops/mac-mini/bin/bill-prediction-collect-scheduled</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__BILL_ROOT__</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__BILL_ROOT__/.rumbling-hedge/logs/launchd-prediction-collect.log</string>
    <key>StandardErrorPath</key>
    <string>__BILL_ROOT__/.rumbling-hedge/logs/launchd-prediction-collect.err.log</string>
  </dict>
</plist>
